# Library/LaunchAgents/self.pg.startup.plist
<!-- For this to work, you need to setenv LC_ALL -->
<!-- 2023-06-01 00:53:25.398 CDT [47945] FATAL:  postmaster became multithreaded during startup -->
<!-- 2023-06-01 00:53:25.398 CDT [47945] HINT:  Set the LC_ALL environment variable to a valid locale. -->
<!-- do this in zshrc: -->
<!-- launchctl setenv LC_ALL "en_us.UTF.8" -->
<!-- why? https://serverfault.com/questions/750827/postgresql-fails-on-startup-fatal-postmaster-became-multithreaded-during-star -->
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>self.pg.startup</string>
    <key>ProgramArguments</key>
    <array>
      <!-- [todo] figure out launchd.conf or something, idk. -->
      <string>/opt/homebrew/opt/postgresql@15/bin/postgres</string>
      <string>-D</string>
      <string>/Users/zooey/Documents/code/site/server/postgres/database</string>      
      <string>-p 9002</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false/>
    </dict>
    <key>StandardErrorPath</key>
    <string>/Users/zooey/.log/launchctl</string>
    <key>StandardOutPath</key>
    <string>/Users/zooey/.log/launchctl</string>
</dict>
</plist>
